# ai-podcast-backend.service
[Unit]
Description=AI Podcast Bot バックエンド (Python)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=waras
WorkingDirectory=/home/waras/ai-podcast/backend
ExecStart=/home/waras/ai-radio/backend/.venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
