# ai-radio-frontend.service
[Unit]
Description=AIラジオ フロントエンド (Bun/Hono)
After=ai-radio-backend.service
Wants=ai-radio-backend.service

[Service]
Type=simple
User=waras
WorkingDirectory=/home/waras/ai-radio/frontend
ExecStart=/home/waras/.bun/bin/bun run src/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=FRONTEND_PORT=3000
Environment=BACKEND_API_URL=http://127.0.0.1:8080

[Install]
WantedBy=multi-user.target
